# init/modprobed-db.timer
[Unit]
Description=Check for new modules

[Timer]
OnUnitActiveSec=6h
