# init/modprobed-db.timer
[Unit]
Description=Check for new modules
PartOf=modprobed-db.service

[Timer]
OnUnitActiveSec=6h
